# deploy/systemd/voltexbazar-web.service
[Unit]
Description=VoltexBazar Next.js frontend
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/apps/voltexbazar
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
StandardOutput=append:/var/log/voltexbazar/web.log
StandardError=append:/var/log/voltexbazar/web-error.log

[Install]
WantedBy=multi-user.target
